# ops/baserow-sync/wlanklar-baserow-sync.timer
[Unit]
Description=Run WLANklar Baserow sync every thirty minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
RandomizedDelaySec=30s
Persistent=true
Unit=wlanklar-baserow-sync.service

[Install]
WantedBy=timers.target
